# dist/dpkg/DEBIAN/postinst
#!/bin/sh

set -ex

deb-systemd-helper enable echo-server.service
